# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=9001

[supervisord]
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisor.pid
environment = KEY1="value1",KEY2="value2"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://localhost:9001

[program:grid_hub]
command=java -jar selenium-server-standalone-2.39.0.jar -role hub
directory=/Users/daniela/selenium
redirect_stderr=true
stdout_logfile=/tmp/grid_hub.log
stderr_logfile=/tmp/grid_hub-error.log

[program:grid_node]
command=java -jar selenium-server-standalone-2.39.0.jar -role node  -nodeConfig config.json -Dwebdriver.chrome.driver="chromedriver"
directory=/Users/daniela/selenium
startsecs=5
redirect_stderr=true
stdout_logfile=/tmp/grid_node.log
stderr_logfile=/tmp/grid_node-error.log
